# src/ghost.conf
# ghost

# description "An Upstart task to make sure that my Ghost server is always running"
# author "Your Name Here"

start on startup

script
    cd /var/www/dustinvtran.com/blog
    npm start --production
end script

